# com.klg-tools.onenote-todo-sync-work.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.klg-tools.onenote-todo-sync-work</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kevin.loyola/Developer/zulu_projects/azure/oneNote_toDo_managment_tool/venv/bin/python3</string>
        <string>src/main.py</string>
        <string>--account</string>
        <string>work</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kevin.loyola/Developer/zulu_projects/azure/oneNote_toDo_managment_tool</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/kevin.loyola/Library/Logs/OneNoteTodoSync/stdout_work.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kevin.loyola/Library/Logs/OneNoteTodoSync/stderr_work.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
